# cron/worker-behavioral.cron
CRON_TZ=Asia/Ho_Chi_Minh
0 * * * * python -m src.worker --job behavioral
