# cron/worker-reports.cron
CRON_TZ=Asia/Ho_Chi_Minh
0 9 * * 1 python -m src.worker --job reports
