# cron/worker-budget.cron
CRON_TZ=Asia/Ho_Chi_Minh
# run hourly at xx:00
0 * * * * python -m src.worker --job budget
