# cron/worker-budget.cron
CRON_TZ=Asia/Ho_Chi_Minh
0 * * * * /usr/local/bin/python -m src.worker --job budget
