# cron/worker-behavioral.cron
CRON_TZ=Asia/Ho_Chi_Minh
# Run daily at 8 AM
0 8 * * * python -m src.worker --job behavioral
